# scripts/com.kline-trainer.daily-work-log.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kline-trainer.daily-work-log</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/lin/dev/KlineTrain/scripts/daily-work-log.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/lin/dev/KlineTrain</string>
    
    <key>StandardOutPath</key>
    <string>/Users/lin/dev/KlineTrain/scripts/logs/daily-work-log.out</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/lin/dev/KlineTrain/scripts/logs/daily-work-log.err</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>